# alpine/mariadb/supervisord.conf
[program:mysqld]
command=/usr/bin/mysqld_safe --user=mysql --console --skip-name-resolve --skip-networking=0
priority=5
