# alpine/mariadb/supervisord.conf
[program:mysqld]
user=mysql
command=/usr/bin/mysqld_safe --console
autostart=true
autorestart=true
priority=5
stdout_events_enabled=true
stderr_events_enabled=true
